# service/offlineimap/run
#!/bin/sh

set -e
set -u

PATH=$PATH:/usr/sbin:/usr/local/bin
HOME=~purzelrakete

exec 2>&1
sleep 15
exec offlineimap
