# service/offlineimap/run
#!/bin/sh

/usr/local/bin/offlineimap 2>&1; sleep 1500

