# service/bitlbee/run
#!/bin/sh

exec /usr/local/sbin/bitlbee -D -n -p 6468 2>&1

