# service/godoc/run
#!/bin/sh

exec godoc -http=:9000 -index=true

