# service/godoc/run
#!/bin/sh

set -e
set -u

PATH=$PATH:/usr/sbin:/usr/local/bin

exec 2>&1
exec godoc -http=:9000 -index=true
